# clients/go/hack/apiban-iptables.timer
[Unit]
Description=APIBan IPTables service schedule

[Timer]
OnActiveSec=300

[Install]
WantedBy=timers.target
