# clients/go/systemd/apiban-iptables.timer
[Unit]
Description=APIBan IPTables service schedule

[Timer]
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
